# deployment/systemd/daphne.service
[Unit]
Description=Lounge Coin Daphne (WebSockets)
After=network.target

[Service]
User=deploy
Group=www-data
WorkingDirectory=/var/www/loungecoin
EnvironmentFile=/var/www/loungecoin/.env
ExecStart=/var/www/loungecoin/venv/bin/daphne \
    -b 127.0.0.1 \
    -p 8001 \
    config.asgi:application
Restart=always
RestartSec=3
Environment="DJANGO_SETTINGS_MODULE=config.settings.production"

[Install]
WantedBy=multi-user.target
